# docker/back/supervisord.conf
[supervisord]
nodaemon=true

[program:fpm]
command=/usr/sbin/php-fpm7.0 -F -c /etc/php/7.0/fpm/php.ini
autorestart=yes
